# installer/debian/DEBIAN/postinst
#!/bin/bash
set -e

if command -v systemctl >/dev/null 2>&1; then
  systemctl daemon-reload
  systemctl enable enigma-agent.service
  systemctl start enigma-agent.service
fi

exit 0