# installer/debian/DEBIAN/prerm
#!/bin/bash
set -e

if command -v systemctl >/dev/null 2>&1; then
  systemctl stop enigma-sensor.service || true
  systemctl disable enigma-sensor.service || true
  systemctl daemon-reload
fi

exit 0